# autostart/prime-offload.desktop
[Desktop Entry]
Type=Application
Exec=prime-offload
Name=Prime Offload
Comment=Required for Optimus Prime Graphics switching
